# warden3/templates/warden_ca_http.service.erb
[Unit]
Description=warden_ca

[Service]
ExecStart=/usr/bin/python /opt/warden_ca/warden_ca_http.py
WorkingDirectory=/opt/warden_ca
Restart=always
# Output to syslog
StandardOutput=syslog
StandardError=syslog
SyslogIdentifier=warden_ca
User=<%= @ca_user %>
Group=<%= @ca_user %>

[Install]
WantedBy=multi-user.target
