# hpdio/templates/dionaea.service.erb
[Unit]
Description=Dionaea Honeypot
Requires=network.target
After=network.target

[Service]
ExecStart=<%= @install_dir %>/bin/dionaea -u <%= @service_user %> -g <%= @service_user %> -w <%= @install_dir %> -p <%= @install_dir %>/var/dionaea.pid

[Install]
WantedBy=multi-user.target
