# hpdio/templates/dionaea.service.erb
[Unit]
Description=Dionaea Honeypot
Requires=network.target
Wants=ntpd.service
After=network.target ntpd.service

[Service]
ExecStart=<%= @install_dir %>/bin/dionaea -u <%= @dio_user %> -g <%= @dio_user %> -w <%= @install_dir %> -p <%= @install_dir %>/var/dionaea.pid

