# hpdio/templates/warden_dio.cron.erb
*/5 * * * * <%= @dio_user %> cd <%= @install_dir %>/warden && ./warden_sender_dio.py 1>/dev/null 2>/dev/null
