# hptelnetd/templates/telnetd.service.erb
[Unit]
Description=telnetd

[Service]
ExecStart=/usr/bin/python <%= @install_dir %>/bin/telnetd.py
ExecStartPost=/usr/bin/sudo <%= @install_dir %>/bin/iptables start
ExecStopPost=/usr/bin/sudo <%= @install_dir %>/bin/iptables stop
WorkingDirectory=<%= @install_dir %>/bin
Restart=always
StandardOutput=syslog
StandardError=syslog
SyslogIdentifier=telnetd.py
User=<%= @service_user %>
Group=<%= @service_user %>

[Install]
WantedBy=multi-user.target
