# hptelnetd/templates/telnetd.service.erb
[Unit]
Description=telnetd

[Service]
ExecStart=/usr/bin/python <%= @install_dir %>/telnetd.py
ExecStartPost=/usr/bin/sudo /opt/telnetd-iptables start
ExecStopPost=/usr/bin/sudo /opt/telnetd-iptables stop
WorkingDirectory=<%= @install_dir %>
Restart=always
StandardOutput=syslog
StandardError=syslog
SyslogIdentifier=telnetd.py
User=<%= @telnetd_user %>
Group=<%= @telnetd_user %>

[Install]
WantedBy=multi-user.target
