# hptelnetd/templates/warden_telnetd.cron.erb
*/5 * * * * <%= @service_user %> cd <%= @install_dir %>/bin && ./warden_sender_telnetd.py 1>/dev/null 2>/dev/null
