# hptelnetd/templates/warden_telnetd.cron.erb
*/5 * * * * <%= @service_user %> cd <%= @install_dir %>/bin; /usr/bin/python warden_sender_telnetd.py 1>><%= @install_dir %>/log/warden_sender_telnetd.log 2>&1
