# hpasa/templates/asa.service.erb
[Unit]
Description=asa

[Service]
ExecStart=/bin/bash -c "/usr/bin/python3 <%= @install_dir %>/bin/asa_server.py --port 443 --enable_ssl --cert <%= @install_dir %>/asacert/bundle.pem &> <%= @install_dir %>/log/asa.log"
WorkingDirectory=<%= @install_dir %>/bin
Restart=always
StandardOutput=syslog
StandardError=syslog
SyslogIdentifier=asa_server.py
User=<%= @service_user %>
Group=<%= @service_user %>

[Install]
WantedBy=multi-user.target
