# hpasa/templates/warden_asa.cron.erb
*/5 * * * * <%= @service_user %> cd <%= @install_dir %>/warden && ./warden_sender_asa.py 1>/dev/null 2>/dev/null
