# hpjdwpd/templates/jdwpd.service.erb
[Unit]
Description=jdwpd

[Service]
ExecStart=/usr/bin/python <%= @install_dir %>/bin/jdwpd.py
WorkingDirectory=<%= @install_dir %>/bin
Restart=always
StandardOutput=syslog
StandardError=syslog
SyslogIdentifier=jdwpd.py
User=<%= @service_user %>
Group=<%= @service_user %>

[Install]
WantedBy=multi-user.target
