# hpucho/templates/uchoudp.service.erb
[Unit]
Description=uchoudp

[Service]
ExecStart=/usr/bin/python <%= @install_dir %>/bin/uchoudp.py
WorkingDirectory=<%= @install_dir %>/bin
Restart=always
StandardOutput=syslog
StandardError=syslog
SyslogIdentifier=uchoudp.py
User=<%= @service_user %>
Group=<%= @service_user %>

Type=idle
LimitNOFILE=65535

[Install]
WantedBy=multi-user.target
