# hpucho/templates/uchotcp.service.erb
[Unit]
Description=uchotcp

[Service]
ExecStart=/usr/bin/python <%= @install_dir %>/uchotcp.py
WorkingDirectory=<%= @install_dir %>
Restart=always
StandardOutput=syslog
StandardError=syslog
SyslogIdentifier=uchotcp.py
User=<%= @uchotcp_user %>
Group=<%= @uchotcp_user %>

Type=idle
LimitNOFILE=65535

[Install]
WantedBy=multi-user.target
