# hpucho/templates/uchoudp.service.erb
[Unit]
Description=uchoudp

[Service]
ExecStart=/usr/bin/python <%= @install_dir %>/uchoudp.py
WorkingDirectory=<%= @install_dir %>
Restart=always
StandardOutput=syslog
StandardError=syslog
SyslogIdentifier=uchoudp.py
User=<%= @uchoudp_user %>
Group=<%= @uchoudp_user %>

Type=idle
LimitNOFILE=65535

[Install]
WantedBy=multi-user.target
