# hpucho/templates/warden_uchotcp.cron.erb
*/5 * * * * <%= @service_user %> cd <%= @install_dir %>/bin; /usr/bin/python warden_sender_uchotcp.py 1>><%= @install_dir %>/log/warden_sender_uchotcp.log 2>&1
