# hpucho/templates/warden_uchoudp.cron.erb
*/5 * * * * <%= @service_user %> cd <%= @install_dir %>/bin; /usr/bin/python warden_sender_uchoudp.py 1>><%= @install_dir %>/log/warden_sender_uchoudp.log 2>&1
