# hpcowrie/templates/cowrie-iptables.service.erb
[Unit]
Description=iptables Cowrie SSH and Telnet Honeypot

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=/bin/sh <%= @install_dir %>/bin/iptables start
ExecStop=/bin/sh <%= @install_dir %>/bin/iptables stop
